# docker/supervisor.conf
[supervisord]
nodaemon=true

[program:app]
command=node /var/www/current/index.js
environment=PORT="3000"
stdout_logfile=/logs/app.log
stdout_logfile_maxbytes=1GB
redirect_stderr=true